# etc/cron.d/logrotate
00 00 * * * root /usr/sbin/logrotate -s /data/logrotate.state /etc/logrotate.conf
